# .config/autostart/keyboard.desktop
[Desktop Entry]
Exec=/home/marco/scripts/keyboard
Icon=dialog-scripts
Name=keyboard
Path=
Type=Application
X-KDE-AutostartScript=true
