# .config/autostart/touchegg.desktop
[Desktop Entry]
Exec=/home/marco/scripts/daemons/touchegg
Icon=dialog-scripts
Name=touchegg
Path=
Type=Application
X-KDE-AutostartScript=true
